# rc.local
# File: /etc/rc.local                            

#
# By default this script does nothing.

# Print the IP address
_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  printf "My IP address is %s\n" "$_IP"
fi

(sleep 10;python /home/pi/pumpcontrol/level-gauge_v3.0b2.py)&
/usr/local/bin/uwsgi --ini /home/pi/projects/django-rpi/pumpcontrol/pumpcontrol$

exit 0
